# docker/root/etc/s6/kmup/run
#!/bin/bash
[[ -f ./setup ]] && source ./setup

pushd /app/kumose/kmup >/dev/null
exec su-exec $USER /usr/local/bin/kmup web
popd
